# crontabs/scripts.sh
#!/bin/bash
php /kunden/homepages/16/d380803416/htdocs/crontabs/testmail.php


php testmail.php
