# crontabs/test.sh
cp testmail.php /kunden/homepages/16/d380803416/htdocs/final/crontabs/backup/
#!/bin/sh
#php /kunden/homepages/16/d380803416/htdocs/crontabs/testmail.php
#php testmail.php